# debian/postinst
#!/bin/sh

set -e

# postinst script for cloudflare-ip-updater

case "$1" in
    configure)
        # Create config file if it doesn't exist
        if [ ! -f /etc/cloudflare-ip-updater/config ]; then
            cp /etc/cloudflare-ip-updater/config.default /etc/cloudflare-ip-updater/config
            chmod 600 /etc/cloudflare-ip-updater/config
            echo "Configuration file created at /etc/cloudflare-ip-updater/config"
            echo "Please edit it with your Cloudflare API token before enabling the service."
        fi

        # Set proper permissions on data directory
        if [ -d /var/lib/cloudflare-ip-updater ]; then
            chown root:root /var/lib/cloudflare-ip-updater
            chmod 755 /var/lib/cloudflare-ip-updater
        fi

        # Reload systemd if needed
        if [ -d /run/systemd/system ]; then
            systemctl daemon-reload >/dev/null || true
        fi

        echo ""
        echo "Cloudflare IP Updater has been installed."
        echo ""
        echo "Next steps:"
        echo "1. Edit /etc/cloudflare-ip-updater/config with your Cloudflare API token"
        echo "2. Enable the service and timer:"
        echo "   sudo systemctl enable cloudflare-ip-updater.service"
        echo "   sudo systemctl enable --now cloudflare-ip-updater.timer"
        echo "3. Check status:"
        echo "   sudo systemctl status cloudflare-ip-updater.timer"
        echo ""
    ;;
esac

exit 0
